# uberspace/app.ini
[program:serlo-app]
directory=~/serlo-editor-as-lti-tool
command=yarn start:dev
autostart=true
autorestart=true
environment=NODE_ENV=production
