# uberspace/app.ini
[program:serlo-app]
directory=/home/%(ENV_USER)s/serlo-editor-as-lti-tool
command=yarn start
autostart=true
autorestart=true
environment=NODE_ENV=production
